# lib/systemd/system/anon-ws-disable-stacked-tor_autogen__var_run_anon-ws-disable-stacked-tor_127.0.0.1_9150.sock.socket
## Copyright (C) 2018 - 2018 ENCRYPTED SUPPORT LP <adrelanos@riseup.net>
## See the file COPYING for copying conditions.

## This file get auto generated by usr/share/anon-ws-disable-stacked-tor/systemd-unit-files-generator.

[Unit]
Description=redirect /var/run/anon-ws-disable-stacked-tor/127.0.0.1_9150.sock to Whonix-Gateway port 9150
Documentation=https://www.whonix.org/wiki/Dev/anon-ws-disable-stacked-tor
ConditionPathExists=!/var/run/qubes/this-is-templatevm

[Socket]
ListenStream=/var/run/anon-ws-disable-stacked-tor/127.0.0.1_9150.sock
SocketUser=debian-tor
SocketMode=666

[Install]
WantedBy=sockets.target
